# templates/le-dns.timer.j2
[Unit]
Description=Refresh letsencrypt certificates regularly

[Timer]
OnCalendar=*-*-* {{ ledns_timing }}

[Install]
WantedBy=timers.target
